# Linux/Backup/wsl_directory/templates/borg-backup.timer
[Unit]
Description=Borg Backup
Requires=borg-backup.service

[Timer]
OnUnitActiveSec=300
Unit=borg-backup.service

[Install]
WantedBy=timers.target
